# tools/ollama.service
[Unit]
Description=Ollama AI Model Server
After=network.target

[Service]
Type=simple
User=%USER%
Environment=OLLAMA_HOST=0.0.0.0:11434
ExecStart=/usr/local/bin/ollama serve
Restart=always
RestartSec=3

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=%h/.ollama
ProtectHome=no

[Install]
WantedBy=multi-user.target</content>
<parameter name="filePath">/Users/fuaadabdullah/ForgeMonorepo/tools/ollama.service
